# systemd/vgkeydesk@.service
[Unit]
Description = VPNGen Keydesk

[Service]
User=%i
Group=%i
WorkingDirectory=/home/%i
AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecStart = /opt/vgkeydesk/keydesk

[Install]
WantedBy=multi-user.target
